# deploy/upstart.conf
description "NGSDB03 Myler Lab NGS application"
start on runlevel [2345]
stop on runlevel [06]
respawn
respawn limit 10 5
exec %(proj_path)s/%(proj_name)s.run_via_gunicorn.sh
